# taiwan_short_term_trading/scripts/com.kelvin.taiwan-paper-pipeline-retry.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kelvin.taiwan-paper-pipeline-retry</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/kelvinhsueh/Documents/taiwan_trading/taiwan_short_term_trading/scripts/run_daily_paper_pipeline.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kelvinhsueh/Documents/taiwan_trading/taiwan_short_term_trading</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>4</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
    <key>PYTHONDONTWRITEBYTECODE</key>
    <string>1</string>
    <key>TAIWAN_PAPER_TRADING_ONLY</key>
    <string>1</string>
    <key>TZ</key>
    <string>America/Chicago</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/kelvinhsueh/Documents/taiwan_trading/taiwan_short_term_trading/logs/launchd_retry_stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kelvinhsueh/Documents/taiwan_trading/taiwan_short_term_trading/logs/launchd_retry_stderr.log</string>
</dict>
</plist>

